# go-media-service/deployments/systemd/go-media-service.service
[Unit]
Description=Campus Pulse Go Media Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/campus-pulse/go-media-service
ExecStart=/opt/campus-pulse/go-media-service/media-service
Restart=always
RestartSec=5
Environment=MEDIA_SERVER_PORT=8090
Environment=MEDIA_PUBLIC_BASE_URL=https://media.example.com
Environment=MEDIA_DATABASE_DSN=./data/media.db
Environment=MEDIA_ADMIN_JWT_SECRET=change-me-admin-secret
Environment=MEDIA_UPLOAD_JWT_SECRET=change-me-upload-secret

[Install]
WantedBy=multi-user.target
