# zens-ops-service/deploy/zens-ops-health.timer
[Unit]
Description=Periodic Zens operations readiness check

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
